# .config/systemd/user/restic-backup.service
[Unit]
Description=Restic backup service for Joe

[Service]
ExecStart=/home/jharry/bin/Backup.sh

[Install]
WantedBy=default.target
