# systemd/.local/share/systemd/user/joetop.mount
[Unit]
Description=Mount Joetop

[Mount]
Where=/mnt/foo
What=/home/user/foo
Type=ext4
