# systemd/.local/share/systemd/user/mnt-worktop-laraplex.automount
[Unit]
Description=mnt-worktop-laraplex

[Automount]
Where=/mnt/worktop/laraplex

[Install]
WantedBy=multi-user.target