# systemd/.local/share/systemd/user/mnt-worktop-interplex.automount
[Unit]
Description=mnt-worktop-interplex

[Automount]
Where=/mnt/worktop/interplex

[Install]
WantedBy=multi-user.target