# systemd/.local/share/systemd/user/mnt-worktop-home.automount
[Unit]
Description=mnt-worktop-home

[Automount]
Where=/mnt/worktop/home

[Install]
WantedBy=default.target