# systemd/.local/share/systemd/user/restic_full.timer
[Unit]
Description=Run Full Restic backup

[Timer]
OnBootSec=20s
Unit=restic_full.service

[Install]
WantedBy=timers.target