# systemd/.local/share/systemd/user/restic.timer
[Unit]
Description=Run Restic backup

[Timer]
OnBootSec=120s
OnCalendar=00/1:00
Unit=restic.service

[Install]
WantedBy=timers.target